# debian/postinst
#!/bin/bash
set -e

systemctl daemon-reload
systemctl enable harboros-subscreen.service
systemctl start harboros-subscreen.service

exit 0
